# pirogue-eve-collector/debian/pirogue-eve-collector.socket
[Unit]
Description=PiRogue socket for Suricata alerts i/o

[Socket]
ListenDatagram=/var/run/suricata.socket

[Install]
WantedBy=sockets.target
