# pirogue-external-exposure/debian/pirogue-external-exposure.timer
[Unit]
Description=Run PiRogue installed certificates renewal

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
