# pirogue-external-exposure/debian/pirogue-external-exposure.service
[Unit]
Description=Renew installed certificates

[Service]
ExecStart=/usr/sbin/renew-certificates
Type=oneshot