# pirogue-external-exposure/debian/pirogue-renew-certificates.service
[Unit]
Description=Renew installed certificates

[Service]
ExecStart=/usr/sbin/renew-certificate
Type=oneshot