# pirogue-maintenance/debian/pirogue-maintenance.service
[Unit]
Description=Run PiRogue maintenance
Wants=suricata.service

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/share/pirogue/pirogue-maintenance.sh
ExecStop=
